# ops/systemd/access-server.service
[Unit]
Description=Claw Access Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/claw/claw-token/server
ExecStart=/usr/bin/node access-server.js
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target